# system/envaya-cleanup.timer
[Unit]
Description=Delete envaya storage daily

[Timer]
OnCalendar=daily
Persistent=true     
 
[Install]
WantedBy=timers.target
